# root/etc/services.d/zoneminder/run
#!/usr/bin/with-contenv bash

eval /usr/bin/zmpkg.pl start

exec \
	s6-setuidgid abc tail -F -q --pid="$(cat /var/run/zm/zm.pid)" \
	/data/log/zmdc.log
